# systemd/heimdall.service
[Unit]
Description=Heimdall Control Plane (nix develop -c ./start.sh)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/blazzee/Development/heimdall_py
EnvironmentFile=/home/blazzee/Development/heimdall_py/.env
ExecStart=/usr/bin/env bash -lc 'nix develop -c ./start.sh'
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
