# systemd/heimdall-agent.service
[Unit]
Description=Heimdall Agent (./start-agent.sh)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/heimdall_py
EnvironmentFile=%h/heimdall_py/.env
Environment=PATH=/usr/bin:/bin:/usr/sbin:/sbin
ExecStart=/bin/bash -lc "./start-agent.sh"
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
